# 03-systemd-service/heartbeat.service
[Unit]
Description=Heartbeat - my first systemd service
After=network.target

[Service]
# Change this path to wherever you cloned the repo
ExecStart=/usr/bin/python3 ~/Linux-Server-Workshop/03-systemd-service/heartbeat.py
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target
